# data/supervisor/laravel_queue.conf
[program:laravel_queue]
command=php /var/www/test/artisan queue:work
numprocs=1
autostart=true
autorestart=true
stderr_logfile=/var/log/laravel_queue.err.log
stdout_logfile=/var/log/laravel_queue.out.log